# cookbooks/zabbix/templates/default/unmonitor_zabbix_host.conf.erb
# Generated by Chef.

description     "Unmonitor the Zabbix host corresponding to this node."
author          "Dhruv Bansal <dhruv@infochimps.com>"

# Trigger on shutdown
start on starting rc RUNLEVEL=[06]

# And it's just a task
task

script
  <%= node.zabbix.conf_dir %>/externalscripts/unmonitor_zabbix_host.rb
  true
end script
